# services/glava.service
[Unit]
Description=Glava
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/glava --desktop

[Install]
WantedBy=default.target
